# packages/contrib/quicksand-alpine/quicksand_alpine/docker/quicksand-guest-agent.ini
[program:quicksand-guest-agent]
command=/usr/local/bin/quicksand-guest-agent
autorestart=true
startretries=9999
stdout_logfile=/var/log/quicksand-guest-agent.log
stderr_logfile=/var/log/quicksand-guest-agent.log
